# systemd/atu-web.service
[Unit]
Description=ATU-100 Web Controller v2.0
After=network.target rigctld.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/atu-controller
ExecStartPre=/home/pi/atu-controller/atu_gpio.py init
ExecStart=/usr/bin/node /home/pi/atu-controller/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
